# home/melden/crontab
# Crontab for *root* user
# Start GPU fan speed controller during start-up
@reboot /home/melden/gpu_speed_fan_control.sh